# vanilla.service
[Unit]
Description=Minecraft Vanilla Server
After=network.target

# which screen || sudo apt-get update && sudo apt-get install -y screen

[Service]
Type=simple
User=connor
WorkingDirectory=/home/connor/MinecraftVanilla
# Start a detached screen session named "mc" running your server
ExecStart=/usr/bin/screen -DmS mc /usr/bin/java -Xmx32G -Xms32G -jar server.jar nogui

# Graceful stop: send "stop<ENTER>" into the screen session
ExecStop=/usr/bin/screen -S mc -p 0 -X stuff "stop\015"

# Give the server time to shut down after "stop"
TimeoutStopSec=60
KillMode=none
Restart=on-failure

[Install]
WantedBy=multi-user.target